# master-server/systemd/kiro-master.service
[Unit]
Description=Kiro Firewall Master Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kiro-master
Group=kiro-master
EnvironmentFile=/etc/kiro-master/master.env
ExecStart=/usr/local/bin/kiro-master
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/kiro-master

[Install]
WantedBy=multi-user.target
