# deployments/systemd/kiro-client-waf.service
[Unit]
Description=Kiro Client WAF - Reverse proxy, challenge pages, rate limiting, and XDP sync
Documentation=https://firewall.vpsgen.com
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/kiro-client-waf --config /etc/kiro/kiro.yaml
WorkingDirectory=/var/lib/kiro

# Restart policy
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=300
StartLimitBurst=5

# User and group (uncomment if 'kiro' user exists on the system)
# User=kiro
# Group=kiro

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
LockPersonality=true

# File system access - needs write for blocklist and logs
ReadWritePaths=/var/lib/kiro /var/log/kiro
ReadOnlyPaths=/etc/kiro /usr/lib/kiro

# Network access required for heartbeat and update check to Master_Server
# Capabilities for XDP/BPF map operations and network binding
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW CAP_NET_BIND_SERVICE

# Resource limits
LimitNOFILE=1048576
LimitNPROC=512

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kiro-client-waf

[Install]
WantedBy=multi-user.target
